# 00-cleanup.conf
[program:container_cleanup]
command=/usr/local/sbin/container_cleanup.sh
priority=-1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=false
startsecs=0
